# deploy/systemd/mini-redis.service
# Install:
#   sudo useradd --system --no-create-home --shell /usr/sbin/nologin mini-redis
#   sudo install -m755 mini-redis /usr/local/bin/mini-redis
#   sudo cp deploy/systemd/mini-redis.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now mini-redis
# Reach it (no public port): ssh -L 6380:localhost:6380 you@vps  →  redis-cli -p 6380

[Unit]
Description=mini-redis-go server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mini-redis
Group=mini-redis
# Bind localhost only — the server has no AUTH yet, so 6380 must never be public.
# Connect over an SSH tunnel. Widen --bind only after adding auth.
ExecStart=/usr/local/bin/mini-redis --bind 127.0.0.1 --port 6380 \
  --aof-path /var/lib/mini-redis/appendonly.aof --appendfsync everysec
Restart=on-failure
RestartSec=1

# systemd creates and owns /var/lib/mini-redis (the AOF lives here).
StateDirectory=mini-redis

# Sandboxing: run with the least privilege the server needs.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
